# config/crontab
# m h  dom mon dow   command
41 6,15 * * *       /usr/bin/update-certs 2>&1
